# network-wait.service
[Unit]
Description=Network Wait Service - Blocks until network is fully operational
Documentation=man:systemd.service(5)
After=sysinit.target systemd-networkd.service NetworkManager.service systemd-resolved.service
Before=network-online.target
Wants=network-online.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/network_monitor.sh --blocking
TimeoutStartSec=900
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=network-online.target